# docker-alertd.service.example
[Unit]
Description=Docker-Alertd
After=docker.service

[Service]
Type=simple
ExecStart=/path/docker-alertd -f /path/docker-alertd.conf.json

[Install]
WantedBy=multi-user.target